# rb24feed/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:rb24feed]
command=/usr/bin/rbfeeder
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
